# sparse/usr/lib/systemd/system/autologin@.service.d/50-vendor-hwcomposer-2-4.conf
[Service]
# Restart hwcomposer to allow user switching to work
ExecStopPost=+-/system/bin/stop vendor.hwcomposer-2-4
ExecStopPost=+-/system/bin/start vendor.hwcomposer-2-4
