# sparse/etc/systemd/system/systemd-ask-password-gui-stop.service.d/50-vendor-hwcomposer-2-4.conf
[Service]
ExecStartPre=/usr/bin/killall systemd-ask-password-gui
# stop first in case something else managed to start it
ExecStartPre=/system/bin/stop vendor.hwcomposer-2-4
ExecStartPre=/system/bin/start vendor.hwcomposer-2-4
ExecStart=/bin/echo After vendor.hwcomposer-2-4 restart
